# scripts/launchd/com.nitish.stocks.fintwitter-finds.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nitish.stocks.fintwitter-finds</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/nitish/stocks automation/scripts/daily_fintwitter_cron.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>15</integer></dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/nitish/stocks automation/data/logs/launchd_fintwitter.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/nitish/stocks automation/data/logs/launchd_fintwitter.err</string>
</dict>
</plist>